# brewhome/brewhome.service
[Unit]
Description=BrewHome - Gestion Brasserie
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/brewhome
ExecStart=/opt/brewhome/venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
